# daily-report/systemd/propnet-daily-report.service
[Unit]
Description=PropNet Daily Report Generator
After=postgresql.service network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/webapp/goldenrabbit/backend/daily-report
ExecStart=/home/webapp/goldenrabbit/backend/venv/bin/python daily_report.py --mode daily
EnvironmentFile=/home/webapp/goldenrabbit/backend/.env
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
